# systemd/oracle.service
[Unit]
Description=Oracle daily job

[Service]
Type=simple
ExecStart=%h/bin/daily.sh

[Install]
WantedBy=default.target
